# net-dmz/web_server/_system/etc/rc.local
#!/bin/bash
/usr/bin/set_default_gw.sh 198.18.1.2
route add -net 198.18.1.128 netmask 255.255.255.128 gw 198.18.1.INNER_DMZ
echo "nameserver 198.18.1.3" > /etc/resolv.conf
route del -host 172.17.0.1
